# etc/services.d/dbus/run
#!/bin/sh

echo STARTING DBUS DAEMON
exec /usr/bin/dbus-daemon --system --nofork --nopidfile --nosyslog
