# etc/services.d/bluetoothd/run
#!/bin/sh

echo STARTING BLUETOOTH DAEMON
exec /usr/lib/bluetooth/bluetoothd --nodetach --noplugin=sap
